# deploy/local/com.mandatehub.operator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- launchd template for the mandatehub operator. Edit CHANGEME values first, then:
       cp com.mandatehub.operator.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.mandatehub.operator.plist
     Testnet default shown; for mainnet set MANDATEHUB_NETWORK=base + the CDP facilitator
     + a CDP key file (gated on H1-H3 for real value). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.mandatehub.operator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/CHANGEME/mandatehub/.venv/bin/python</string>
    <string>/Users/CHANGEME/mandatehub/deploy/local/operator.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>MANDATEHUB_FACILITATOR_URL</key><string>https://x402.org/facilitator</string>
    <key>MANDATEHUB_NETWORK</key><string>base-sepolia</string>
    <key>MANDATEHUB_PAY_TO</key><string>0xCHANGEME</string>
    <key>MANDATEHUB_PORT</key><string>8403</string>
    <key>MANDATEHUB_PUBLIC_URL</key><string>https://mandatehub.CHANGEME</string>
    <!-- mainnet: MANDATEHUB_FACILITATOR_URL=https://api.cdp.coinbase.com/platform/v2/x402
                  MANDATEHUB_NETWORK=base
                  MANDATEHUB_CDP_KEY_FILE=/Users/CHANGEME/.mandatehub-cdp.json -->
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/CHANGEME/.mandatehub-operator/operator.log</string>
  <key>StandardErrorPath</key><string>/Users/CHANGEME/.mandatehub-operator/operator.log</string>
</dict>
</plist>
