# deploy/local/com.mandatehub.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- mandatehub backup — edit CHANGEME, cp to ~/Library/LaunchAgents/, launchctl load it. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.mandatehub.backup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/CHANGEME/mandatehub/.venv/bin/python</string>
    <string>/Users/CHANGEME/mandatehub/deploy/local/backup.py</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/CHANGEME/.mandatehub-operator/backup.log</string>
  <key>StandardErrorPath</key><string>/Users/CHANGEME/.mandatehub-operator/backup.log</string>
</dict>
</plist>
